# deploy/systemd/gh-social-ml.service
[Unit]
Description=GH Social ML API Service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f gh-social-ml
ExecStart=/usr/bin/docker run --rm --name gh-social-ml --network host --env-file /etc/gh-social/ml.env gh-social-ml:current uvicorn api.main:app --host 127.0.0.1 --port 8000
ExecStop=/usr/bin/docker stop --time 30 gh-social-ml
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
